# appconfig/dss.service
[Unit]
Description=DSS service

[Service]
Type=forking
User=dataiku
ExecStart=/home/dataiku/dss/bin/dss start

[Install]
WantedBy=multi-user.target
